# root/etc/profile.d/umask.sh
# previously set by /etc/profile
umask 022
