# root/etc/systemd/system/getty@.service.d/override.conf
[Service]
ExecStart=
# '%I' was replaced with '-' in	/usr/lib/systemd/system/getty@.service with systemd 250-4
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear %I $TERM

