# config/systemd/user/rclone-sync.service
[Unit]
Description=rclone sync

[Service]
Type=oneshot
ExecStart=/home/alex/bin/rclone-sync

