# config/systemd/user/rclone-sync.timer
[Unit]
Description=rclone sync timer

[Timer]
OnBootSec=5
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
