# etc/profile.d/homebin.sh
#!/bin/sh

if [ -d ~/bin ]; then
	export PATH=~/bin:${PATH}
fi
