# root/etc/systemd/user/clean-home.service
[Unit]
Description=clean home

[Service]
Type=oneshot
ExecStart=/home/alex/bin/clean-home

