# root/etc/systemd/user/rival-battery.service
[Unit]
Description=rival battery

[Service]
Type=oneshot
ExecStart=/home/alex/bin/rival-battery

