# config/autostart/pasystray.desktop
[Desktop Entry]
Version=1.0
Name=PulseAudio System Tray
GenericName=
Comment=An Applet for PulseAudio
Exec=pasystray --volume-max=100 --volume-inc=5 --no-icon-tooltip --notify=none
Icon=pasystray
StartupNotify=true
Type=Application
Categories=AudioVideo;Audio;
Keywords=pulseaudio;tray;system tray;applet;volume;
