# root/etc/systemd/system/cpu-power.service
[Unit]
Description=turbostat CPU power

[Service]
ExecStart=/home/alex/bin/cpu-power
Restart=no

[Install]
WantedBy=multi-user.target
